# supervisord.conf
[supervisord]
nodaemon=true

[program:squid3]
command=/apps/squid/sbin/squid  -NsY
stdout_events_enabled=true
stderr_events_enabled=true


